# init.d/xray
#!/sbin/openrc-run

name="Xray Daemon"
description="A unified platform for anti-censorship"

supervisor=supervise-daemon
respawn_delay=2
respawn_max=0

pidfile="/run/${RC_SVCNAME}.pid"
rc_ulimit="-n 1024000 -u 1024000"
capabilities="cap_net_bind_service,cap_net_admin,cap_net_raw"
extra_commands="checkconfig"

env=${env:-"XRAY_LOCATION_ASSET=/usr/local/share/xray/"} # See: https://github.com/XTLS/Xray-core/blob/main/common/platform/platform.go
confdir=${confdir:-"/usr/local/etc/xray/"}
command="/usr/local/bin/xray"
command_user="nobody:nobody"
command_args="run -confdir $confdir"
required_files="$confdir"

depend() {
    need net
    want dns ntp-client
    after firewall
}

checkconfig() {
    export $env
    $command $command_args -test
    eend $?
}

start_pre() {
    checkconfig || return 1
}
