# init.d/xray
#!/sbin/openrc-run
supervisor=supervise-daemon

name="Xray"
description="The best v2ray-core, with XTLS support"
description_checkconfig="Test configuration file"

: ${env:="XRAY_LOCATION_ASSET=/usr/local/share/xray/"}
: ${confdir:="/usr/local/etc/xray/"}

command="/usr/local/bin/xray"
command_args="run -confdir $confdir"
command_user="nobody"
extra_commands="checkconfig"
capabilities="^cap_net_bind_service"

depend() {
	need net localmount
	after firewall
}

checkconfig() {
	if [ ! -d "$confdir" ]; then
		eerror "You need to setup $confdir first"
		return 1
	fi
	export $env
	su ${command_user%:*} -s /bin/sh -c "$command $command_args -test"
	eend $?
}

stop_pre() {
	if [ "$RC_CMD" = restart ]; then
		checkconfig || return $?
	fi
}
